# systemd/crt-player.service
[Unit]
Description=CRT Player — telnet control server for pi-crt-player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# @USER@ is substituted with the installing user by setup.sh.
User=@USER@
ExecStart=/usr/bin/python3 /usr/local/lib/pi-crt-player/player.py
Environment=PLAYER_PORT=23
# Allow binding the privileged telnet port (23) without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
